# crontab
0 23 * * * root cd /app && python jellyfin-youtube-downloader.py >> /proc/1/fd/1 2>&1

